# escribe-normas.service
[Unit]
Description = Servicio escribe normas

[Service]
ExecStart=/bin/bash /home/student/normas.sh

[Install]
WantedBy=multi-user.target
